# com.herculeanfit1.aistudyplans.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.herculeanfit1.aistudyplans.backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/herculeanfit1/dev/AIStudyPlans/aistudyplans-backup_script.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>21600</integer>
    <key>StandardOutPath</key>
    <string>/Users/herculeanfit1/dev/AIStudyPlans/aistudyplans-launchd_backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/herculeanfit1/dev/AIStudyPlans/aistudyplans-launchd_backup_error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist> 